# shell/profile.d/01-colors.sh
# Define some colors to be used in follinwg scripts

RST="\e[0m"         # Reset
RVS="\e[7m"         # Reverse
RSTRVS="\e[27m"     # Reset Reverse
BLD="\e[1m"         # Bold
GRN="\e[32m"        # Green
WHT="\e[97m"        # White
BLU="\e[34m"        # Blue
YLW="\e[93m"        # Yellow
RED="\e[91m"        # Red
VLT="\e[38;5;55m"   # Violet
BGVLT="\e[48;5;55m" # Background Violet
BGRST="\e[49m"      # Background default
