# shell/profile.d/01-styles.sh
# Define some styles to be used in following scripts

STL_RST="\e[0m"     # Reset
STL_BLD="\e[1m"     # Bold
STL_YLW="\e[93m"    # Yellow
