# shell/profile.d/01-styles.sh
# Define some styles to be used in following scripts
# (only if running in a terminal)
if test -t 0; then
    STL_RST="\e[0m"     # Reset
    STL_BLD="\e[1m"     # Bold
    STL_YLW="\e[93m"    # Yellow
fi
