# shell/profile.d/99-hooks.sh
# Setup hooks for prompt and pre-command functions

# Hook function to reckon PS1 every time
PROMPT_COMMAND=__updateMyPrompt
