# shell/profile.d/work.sh
# Some work-specific settings
export REMOTEUSER=dpochet
